# deploy/gunicorn.service.template
[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/root/$PROJECT_DIR
ExecStart=gunicorn \
          --access-logfile /root/$PROJECT_DIR/access.log \
          --error-logfile /root/$PROJECT_DIR/error.log \
          --workers 3 \
          --capture-output --log-level info \
          --bind unix:/run/gunicorn.sock \
          $PROJECT_DIR.wsgi:application

[Install]
WantedBy=multi-user.target